# fmberry.service
[Unit]
Description=FMBerry

[Service]
Type=simple
ExecStart=/usr/local/bin/fmberryd nodaemon

[Install]
WantedBy=multi-user.target
